# deploy/systemd/bonus_expiration_warning.timer
[Unit]
Description=Run Referral MVP bonus expiration warnings daily

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
