# deploy/systemd/referral_admin.service
[Unit]
Description=Referral MVP Django Admin
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/sistema_loialnosti/backend
EnvironmentFile=/opt/sistema_loialnosti/.env
ExecStart=/opt/sistema_loialnosti/backend/.venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
